# srcs/nginx/srcs_nginx/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=sh -c "kill $(pidof nginx) || nginx -g 'daemon off;' && kill $(pidof supervisord)"
command=sh -c "nginx -g 'daemon off;' && kill -s SIGTERM $(cat supervisord.pid)"
autostart=true
autorestart=true

[program:sshd]
command=sh -c "kill $(pidof sshd) || /usr/sbin/sshd -D && kill $(pidof supervisord)"
command=sh -c "/usr/sbin/sshd -D && kill $(pidof supervisord)"
autostart=true
autorestart=true
